# cron/epylog.cron.in
#!/bin/sh
# Run epylog daily.
#
# $Id$
#
%%sbindir%%/epylog --cron
